# foris-controller-mock/supervisor/foris-controller.conf
[program:foris-controller]
directory=/root/
command=/usr/local/bin/foris-controller --backend mock mqtt --controller-id 0000000000000000
