# reforis-demo/supervisor/foris-ws.conf
[program:foris-ws]
directory=/var/www/
user=www-data
command=/usr/local/bin/foris-ws -a none --host 127.0.0.1 --port 9081 mqtt
