# reforis-demo/supervisor/reforis.conf
[program:reforis]
directory=/var/www/
environment = FLASK_APP="reforis:create_app('demo')",FLASK_ENV="development",CONTROLLER_ID="0000000000000000"
user=www-data
command=/usr/local/bin/reforis
