# supervised/supervisor/mosquitto.conf
[program:mosquitto]
directory=/var/lib/mosquitto/
command=/usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
